# flatpak-auto-update.timer
# ~/.config/systemd/user/flatpak-auto-update.timer
# Timer unit for daily Flatpak updates
[Unit]
Description=Run daily Flatpak updates

[Timer]
# Run daily
OnCalendar=daily
# Run missed updates after system reboot (if system was down during a scheduled update)
Persistent=true
# Add random delay up to 1 hour to avoid all systems updating at once
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
