# flatpak-auto-update.service
# ~/.config/systemd/user/flatpak-auto-update.service
# Service unit for updating Flatpak apps
[Unit]
Description=Update Flatpak apps
# Wait for network to be available (non-blocking)
Wants=network-online.target
# Only start after network is online
After=network-online.target

[Service]
# Service runs once and exits (not a daemon)
Type=oneshot
# Update Flatpaks --noninteractive: no prompts, -y: answer yes to all
ExecStart=/usr/bin/flatpak update --noninteractive -y
# Prevent the service from gaining additional privileges
NoNewPrivileges=true
